# run-on-reboot.service
[Unit]
Description=Run a script on reboot

[Service]
Type=oneshot
ExecStart=/path/to/your/script.sh

[Install]
WantedBy=multi-user.target
